# packaging/systemd/rsdbd.socket
[Unit]
Description=RSDB TCP socket

[Socket]
ListenStream=27101
NoDelay=true

[Install]
WantedBy=sockets.target
